# root/etc/s6-overlay/s6-rc.d/unbound/finish
#!/command/execlineb -P
killall -9 unbound
