# root/etc/s6-overlay/s6-rc.d/unbound/run
#!/command/execlineb -P
/usr/sbin/unbound -d
